# templates/consul.initd.erb
#! /bin/sh
### BEGIN INIT INFO
# Provides:          consul
# Required-Start:    $all
# Required-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:
# Short-Description: start consul
### END INIT INFO


PATH=/sbin:/usr/sbin:/bin:/usr/bin
CONSUL=<%= scope.lookupvar('consul::bin_dir') -%>/consul
CONFIG=<%= scope.lookupvar('consul::config_dir') -%>


do_start() {
	if [ -x $CONSUL ]; then
		$CONSUL agent -config-dir $CONFIG
		ES=$?
		return $ES
	fi
}

case "$1" in
    start)
	do_start
        ;;
    restart|reload|force-reload)
        echo "Error: argument '$1' not supported" >&2
        exit 3
        ;;
    stop)
        ;;
    *)
        echo "Usage: $0 start|stop" >&2
        exit 3
        ;;
esac
